# crontab
1 0 * * 1-5 ./capture-screenshot.sh http://www.google.com paulin.shek@gmail.com